# usizo0/00-configure-usizo/files/crontab
# call home
@reboot nobody /usr/bin/usizo-watchdog.sh
*/5 * * * * nobody /usr/bin/usizo-watchdog.sh
